# systemd/user/bakop.timer
[Unit]
Description=backup every 30 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=120min

[Install]
WantedBy=timers.target
