# systemd/user/sync-mail.service
[Unit]
Description=Sync and index mail

[Service]
Type=oneshot
ExecStart=/home/engon/.config/bin/sync-mail
