# systemd/user/waybar-restart.service
[Unit]
Description=Restart Waybar
Documentation=https://github.com/Alexays/Waybar/wiki
[Service]
Type=oneshot
ExecStart=/usr/bin/pkill -USR2 waybar
